# crownpipe/media/systemd/bgremove.service
[Unit]
Description=Background removal batch processing
After=network.target

[Service]
Type=oneshot
User=crown-pipeline
Group=marketing
Environment=HOME=/home/crown-pipeline
Environment=PYTHONPATH=/opt/crownpipe
WorkingDirectory=/opt/crownpipe
ExecStart=/opt/crownpipe/venv/bin/python3 /opt/crownpipe/bin/media-bgremove

[Install]
WantedBy=multi-user.target
