# crownpipe/media/systemd/bgremove.timer
[Unit]
Description=Run background removal every minute

[Timer]
OnBootSec=45
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
