# systemd/crownpipe-deploy-production.timer
[Unit]
Description=Run production deployment every 5 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
