# systemd/crownpipe-prepare-formatting.timer
[Unit]
Description=Run prepare formatting every minute

[Timer]
OnBootSec=75
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
